# hactar.conf
# hactar project
start on runlevel 2
stop on runlevel [!2]
respawn
setuid hactar
setgid hactar

exec /home/hactar/git/hactar/run.py
