# hactar.conf
# hactar project
start on runlevel 2
stop on runlevel [!2]
respawn
setuid hactar
setgid hactar

pre-start script
    cd /home/hactar/git/hactar/
end script

exec /home/hactar/git/hactar/run.py
