# etc/hactar.conf
# hactar project
start on runlevel 2
stop on runlevel [!2]
respawn
setuid vagrant
setgid vagrant

chdir /home/vagrant/hactar

exec /home/vagrant/run.py
